# extra/linux-devkit/arm-none-linux-gnueabi/etc/xinetd.d/cups-lpd
service printer
{
	socket_type = stream
	protocol = tcp
	wait = no
	user = lp
	group = lp
	passenv =
	server = /usr/lib/cups/daemon/cups-lpd
	server_args = -o document-format=application/octet-stream
}
